# configs/etc/s6/chronograf/run
#!/bin/bash

/opt/chronograf/chronograf -config=/opt/etc/chronograf.toml
